# deploy/transients.service
[Unit]
Description=transients
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=transients
Group=transients
ExecStart=/usr/local/bin/transients
EnvironmentFile=/etc/transients/env
WorkingDirectory=/var/lib/transients
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/var/lib/transients

[Install]
WantedBy=multi-user.target
